# sonarqube-server/sonarqube-server.service
[Unit]
Description=SonarQube Server (Podman rootless)
Documentation=https://docs.sonarqube.org/latest/
After=network.target

[Service]
Type=simple
Restart=on-failure
RestartSec=5
User=%u
StandardOutput=journal
StandardError=journal

# Load configuration
EnvironmentFile=%h/.config/sonarqube/sonarqube.env

# Create data directory
ExecStartPre=mkdir -p /var/lib/sonarqube
ExecStartPre=chmod 755 /var/lib/sonarqube

# Pull latest image
ExecStartPre=/usr/bin/podman pull ${SONARQUBE_IMAGE:-sonarqube:latest}

# Remove existing container if present
ExecStartPre=-/usr/bin/podman rm -f ${SONARQUBE_CONTAINER_NAME:-sonarqube}

# Start SonarQube container with proper configuration
ExecStart=/usr/bin/podman run \
  --name ${SONARQUBE_CONTAINER_NAME:-sonarqube} \
  --rm \
  --memory ${SONARQUBE_MEMORY_LIMIT:-2g} \
  --cpus ${SONARQUBE_CPUS_LIMIT:-2} \
  -p ${SONARQUBE_PORT:-9000}:9000 \
  -p 9092:9092 \
  -v /var/lib/sonarqube/postgresql:/var/lib/postgresql \
  -v /var/lib/sonarqube/sonarqube:/opt/sonarqube/data \
  -v /var/lib/sonarqube/extensions:/opt/sonarqube/extensions \
  -v /var/lib/sonarqube/logs:/opt/sonarqube/logs \
  -e SONAR_ES_BOOTSTRAP_CHECKS_DISABLED=true \
  -e SONAR_TELEMETRY_ENABLE=false \
  ${SONARQUBE_IMAGE:-sonarqube:latest}

# Stop container cleanly
ExecStop=/usr/bin/podman stop ${SONARQUBE_CONTAINER_NAME:-sonarqube}

[Install]
WantedBy=default.target
